# apps/dashboard/crontab
0 6 * * * node --max-old-space-size=1536 server/daily-leads.js